# src/pmproxy/pmproxy_check.service.in
[Unit]
Description=Optionally compress remote PCP archives from pmproxy
Documentation=man:pmlogger_daily(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Environment="PMLOGGER_DAILY_PARAMS="
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmproxy_timers
# proxy-managed archives only (-Q) and compress (-K)
ExecStart=@PCP_BINADM_DIR@/pmlogger_daily --noerror -Q -K $PMLOGGER_DAILY_PARAMS
WorkingDirectory=@PCP_VAR_DIR@
Group=@PCP_GROUP@
User=@PCP_USER@
## Paths
ProtectProc=invisible
# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=
# Security
NoNewPrivileges=true
# Process properties
UMask=077
# Sandboxing ...
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateMounts=true
# System call filtering
SystemCallFilter=@system-service
# Note: pmproxy needs @aio
SystemCallFilter=~@resources @keyring @memlock
SystemCallArchitectures=native
