# src/pmlogger/pmlogger_farm_check.service.in
[Unit]
Description=Check PCP archive loggers
Documentation=man:pmlogctl(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=@SD_SERVICE_TYPE@
Restart=no
TimeoutStartSec=4h
TimeoutStopSec=120
# needs to run as root:root not pcp:pcp, pmlogctl will handle the
# switch to pcp:pcp if any pmlogger needs to be launched, so NO
# User= nor Group= here
# -m is for migrate that we don't really need anymore
# check will check and if necessary restart all pmlogger instances
# known to pmlogctl(1)
ExecStart=@PCP_BIN_DIR@/pmlogctl -m check
WorkingDirectory=@PCP_VAR_DIR@
## Paths
ProtectProc=invisible
# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=
# Security
NoNewPrivileges=true
# Process properties
UMask=077
# Sandboxing ...
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateMounts=true
# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@aio @resources @keyring @memlock
SystemCallArchitectures=native
