# src/pmlogger/pmlogger.service.in
[Unit]
Description=PCP primary archive logger
Documentation=man:pmlogger(1)
After=network-online.target pmcd.service
Before=pmlogger_farm.service pmlogger_check.timer pmlogger_daily.timer
BindsTo=pmlogger_farm.service pmlogger_check.timer pmlogger_daily.timer
Wants=pmcd.service

[Service]
Type=notify
NotifyAccess=all
# The pmlogger service manages the primary pmlogger (only)
Environment="PMLOGGER_CHECK_PARAMS=--only-primary"
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmlogger
TimeoutStartSec=120
TimeoutStopSec=120
Restart=always
ExecStart=@PCP_SERVICES_DIR@/pmlogger start-systemd
ExecStop=@PCP_SERVICES_DIR@/pmlogger stop-systemd
PIDFile=@PCP_RUN_DIR@/pmlogger.pid
Group=@PCP_GROUP@
User=@PCP_USER@
## Paths
ProtectProc=invisible
# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=
# Security
NoNewPrivileges=true
# Process properties
UMask=077
# Sandboxing ...
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateMounts=true
# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@aio @resources @keyring @memlock
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
